# systemd/tgbot.service
# Example configuration for running the bot on a server as a systemd service
[Unit]
Description=OpenWeatherBot daemon
After=network.target

[Service]
User=tgbot
Group=tgbot
Type=simple
WorkingDirectory=/opt/open-weather-bot
ExecStart=/opt/open-weather-bot/venv/bin/python /opt/open-weather-bot/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
